# root/etc/uci-defaults/40_luci-k-permission
#!/bin/sh

[ -e "/etc/init.d/rpcd" ] && /etc/init.d/rpcd restart

rm -f /tmp/luci-indexcache

exit 0
